# infra/deploy/systemd/esdata-esma-firds-daily.timer
[Unit]
Description=Run esdata ESMA FIRDS daily sync

[Timer]
OnCalendar=*-*-* 04:20:00 Europe/Madrid
Persistent=true
Unit=esdata-job@cron-esma-firds-daily.service

[Install]
WantedBy=timers.target
