# infra/deploy/systemd/esdata-job@.service
[Unit]
Description=esdata scheduled job (%i)
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=deploy
WorkingDirectory=/opt/esdata
ExecStart=/usr/bin/docker compose --env-file /etc/esdata/esdata.env -f /opt/esdata/infra/deploy/docker-compose.prod.yml run --rm %i
